# service/privoxy/run
#!/bin/sh

exec /usr/sbin/privoxy --user privoxy --no-daemon
